# config/treadmill-sensor.service
[Unit]
Description=Treadmill Usage Monitoring Sensor

[Service]
ExecStart=/usr/src/treadmill/treadmill
Type=simple

