# config/treadmill-daily.service
[Unit]
Description=Treadmill

[Service]
RemainAfterExit=yes
ExecStart=/usr/src/treadmill/treadmill-daily
Type=oneshot
